# com.glasshouse.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.glasshouse.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mabramsky/glasshouse1/scripts/daily_run.sh</string>
    </array>

    <!-- Run daily at 12:00 PM (noon) - after manual Parcl data drop at 9am -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>12</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mabramsky/glasshouse1/logs/launchd_daily.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mabramsky/glasshouse1/logs/launchd_daily_error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/mabramsky/glasshouse1</string>

    <!-- Run missed jobs when waking from sleep -->
    <key>StartInterval</key>
    <integer>86400</integer>
</dict>
</plist>
