# com.glasshouse.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.glasshouse.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mabramsky/glasshouse/scripts/weekly_run.sh</string>
    </array>

    <!-- Run every Monday at 7:30 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mabramsky/glasshouse/logs/launchd_weekly.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mabramsky/glasshouse/logs/launchd_weekly_error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/mabramsky/glasshouse</string>
</dict>
</plist>
